# script/lianjia_crawler_job.cron
0 5 * * * ../crawler/lianjia.py >> ../log/logfile