# script/lianjia_crawler_job.cron
0 18 * * * ../crawler/lianjia.py bj >> ../log/logfile